# phase2_application_lockdown/exam-watchdog.service
[Unit]
Description=Exam Process Watchdog
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/exam-watchdog.sh
Restart=always

[Install]
WantedBy=multi-user.target
